# rc.local
#!/bin/bash

[ -f /opt/closure/installed ] || bash /opt/closure/init.sh
